# EEGManager/EEGManager/Docker/supervisord.conf
[supervisord]
nodaemon=true

[program:jupyter]
command=jupyter notebook 
redirect_stderr=true

[program:tensorboard]
command=tensorboard --logdir /content/logs/
redirect_stderr=true